# config/a2host.crontab
0 0 * * * /home/django/scripts/daily_maintence.sh
